# user/systemd/user/archive-backup.timer
[Unit]
Description=Create a backup of essentials every sunday at 10pm

[Timer]
OnCalendar=Sat *-*-* 22:00:00

[Install]
WantedBy=X.target
