# user/systemd/user/archive-backup.service
[Unit]
Description=Create a backup of essentials

[Service]
Type=oneshot
ExecStart=/bin/zsh -c "/home/xylous/bin/barckup"
