# user/systemd/user/wallpaper-change.timer
[Unit]
Description=Change wallpapers every hour

[Timer]
OnActiveSec=0sec
OnUnitActiveSec=60min

[Install]
WantedBy=X.target
