# user/systemd/user/wallpaper-change.timer
[Unit]
Description=Change wallpapers every once in a while

[Timer]
OnActiveSec=0sec
OnUnitActiveSec=45min

[Install]
WantedBy=X.target
