# user/systemd/user/X.target
[Unit]
Description=Xorg server is running
